# Actividad4/actividad4.service
[Unit]
Description=Saludo y fecha

[Service]
Type=simple
ExecStart=/usr/local/bin/actividad4.sh

[Install]
WantedBy=multi-user.target